# ebm.conf
#!upstart
description "TECOR - Microclim Config"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

chdir /home/ubuntu/ebm

env NODE_ENV=production
env PORT=80
env MONGODB_URI=mongodb://localhost:27017/ebm

# exec node /home/ubuntu/ebm/bin/www >> /home/ubuntu/ebm/app.log 2>&1  


